# tripapp.service
[Unit]
Description=TripCost App
After=network.target

[Service]
WorkingDirectory=/opt/tripapp
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
